# root/etc/services.d/sopel/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc sopel --config /config/.sopel/default.cfg